# deploy/systemd/rdashboard-source-ingress-bridge.socket
[Unit]
Description=Private Kamal network listener for rdashboard source ingress
After=network-online.target docker.service
Wants=network-online.target

[Socket]
ListenStream=172.19.0.1:3201
NoDelay=true

[Install]
WantedBy=sockets.target
