# deploy/systemd/rdashboard-notify.service
[Unit]
Description=rdashboard isolated Telegram notification delivery
After=network-online.target local-fs.target
Wants=network-online.target
Before=rdashboard.service

[Service]
Type=simple
User=rdashboard-notify
Group=rdashboard-notify
ExecStart=/var/lib/rdashboard-bootstrap/current/bin/rdashboard-notify
Environment=RDASHBOARD_NOTIFY_DATA_DIR=/var/lib/rdashboard-notify
Environment=RDASHBOARD_NOTIFY_SOCKET=/run/rdashboard-notify/notify.sock
EnvironmentFile=/etc/rdashboard/notifier.env
LoadCredential=telegram-gateway-secret:/etc/rdashboard/credentials/telegram-gateway-secret
RuntimeDirectory=rdashboard-notify
RuntimeDirectoryMode=0750
StateDirectory=rdashboard-notify
StateDirectoryMode=0700
UMask=0007
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
TasksMax=64
LimitNOFILE=1024
MemoryMax=64M
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
