# deploy/systemd/rdashboard-workflow-launcher.service
[Unit]
Description=rdashboard fixed root workflow job launcher
After=local-fs.target dbus.socket
Before=rdashboard-worker.service

[Service]
Type=simple
User=root
Group=rdashboard-worker
ExecStart=/var/lib/rdashboard-bootstrap/current/bin/rdashboard-workflow-launcher
RuntimeDirectory=rdashboard-workflow-launcher
RuntimeDirectoryMode=0750
UMask=0077
NoNewPrivileges=yes
PrivateDevices=yes
PrivateNetwork=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
ReadOnlyPaths=/var/lib/rdashboard-build/titanium
ReadWritePaths=/var/lib/rdashboard-workflow-launcher /var/lib/rdashboard-build/operations /var/lib/rdashboard-build/oci-results /var/lib/rdashboard-build/self-releases /var/lib/rdashboard-build/titanium/roots/active-operation /run/rdashboard-workflow-launcher
InaccessiblePaths=-/etc/rdashboard/credentials -/run/rdashboard -/run/rdashboard-source-ingress -/run/rdashboard-source-delivery -/run/rdashboard-workflow -/run/docker.sock -/var/run/docker.sock -/run/containerd
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
CapabilityBoundingSet=CAP_CHOWN CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER
AmbientCapabilities=CAP_CHOWN CAP_DAC_OVERRIDE CAP_DAC_READ_SEARCH CAP_FOWNER
DevicePolicy=closed
TasksMax=64
LimitNOFILE=256
MemoryMax=256M
CPUQuota=100%
IOWeight=25
TimeoutStartSec=10s
TimeoutStopSec=15s
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
