# deploy/systemd/rdashboard-dependency-fetcher.service
[Unit]
Description=rdashboard fixed public crates.io dependency fetcher
After=network-online.target
Before=rdashboard-worker.service
Wants=network-online.target

[Service]
Type=simple
User=rdashboard-dependency-fetcher
Group=rdashboard-dependency-fetch
ExecStart=/usr/libexec/rdashboard/rdashboard-dependency-fetcher
EnvironmentFile=/etc/rdashboard/workflow-worker.env
RuntimeDirectory=rdashboard-dependency-fetcher
RuntimeDirectoryMode=0750
UMask=0077
NoNewPrivileges=yes
PrivateDevices=yes
PrivateTmp=yes
ProtectClock=yes
ProtectControlGroups=yes
ProtectHome=yes
ProtectHostname=yes
ProtectKernelLogs=yes
ProtectKernelModules=yes
ProtectKernelTunables=yes
ProtectSystem=strict
InaccessiblePaths=-/etc/rdashboard/credentials -/var/lib/rdashboard -/var/lib/rdashboard-build -/var/lib/rdashboard-workflow-launcher -/run/rdashboard -/run/rdashboard-source-ingress -/run/rdashboard-source-delivery -/run/rdashboard-workflow -/run/rdashboard-workflow-launcher -/run/docker.sock -/var/run/docker.sock -/run/containerd -/run/podman/podman.sock
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6
RestrictNamespaces=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SystemCallArchitectures=native
CapabilityBoundingSet=
AmbientCapabilities=
DevicePolicy=closed
UnsetEnvironment=ALL_PROXY HTTP_PROXY HTTPS_PROXY NO_PROXY all_proxy http_proxy https_proxy no_proxy
TasksMax=64
LimitNOFILE=256
MemoryMax=384M
MemorySwapMax=0
CPUQuota=100%
IOWeight=50
TimeoutStartSec=15s
TimeoutStopSec=15s
Restart=on-failure
RestartSec=2s

[Install]
WantedBy=multi-user.target
